# backend/docker/s6/nginx/run
#!/command/execlineb -P
with-contenv
exec nginx -g "daemon off;"
