# backend/docker/s6/setup_groups/up
#!/command/with-contenv sh

cd /app

uv run manage.py setup_groups