# backend/docker/s6/migrate/up
#!/command/execlineb -P
cd /app
exec uv run manage.py migrate --noinput