# backend/docker/s6/gunicorn/run
#!/command/with-contenv sh

cd /app

fdmove -c 2 1 uv run gunicorn -k uvicorn.workers.UvicornWorker -w 2 -b unix:/run/gunicorn.sock mai.asgi:application
